# cnest/data/profile.d/stowhome.sh
if [ -x "$HOME/.stow_profile" ]; then
    "$HOME/.stow_profile"
fi
